# gunicorn.simple.service
[Unit]
Description=GBot Web Application
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/opt/gbot-web-app-original-working
Environment="PATH=/opt/gbot-web-app-original-working/venv/bin"
ExecStart=/opt/gbot-web-app-original-working/venv/bin/gunicorn \
          --workers 3 \
          --bind 127.0.0.1:5000 \
          --error-logfile /var/log/gbot-error.log \
          --access-logfile /var/log/gbot-access.log \
          app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
