# gunicorn.fixed.service
[Unit]
Description=GBot Web Application
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/gbot-web-app-original-working
Environment="PATH=/opt/gbot-web-app-original-working/venv/bin"
# Create logs directory first
ExecStartPre=/bin/mkdir -p /opt/gbot-web-app-original-working/logs
# Set proper permissions for logs directory
ExecStartPre=/bin/chown -R www-data:www-data /opt/gbot-web-app-original-working/logs
# Start gunicorn with direct log output to avoid file not found error
ExecStart=/opt/gbot-web-app-original-working/venv/bin/gunicorn \
          --workers 3 \
          --bind 127.0.0.1:5000 \
          --error-logfile - \
          --access-logfile - \
          app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
